# crontab
*/30 * * * * root cp -r /var/www/owncloud/3rdparty /data/

